# systemd/gamesrv@.service
# /etc/systemd/system/gamesrv@.service
# Template unit: one instance per game server.  Instance name (%i) matches the
# server's `name:` in servers/<name>.yml AND its install dir under /srv/gameservers.
#
# Enable:  systemctl enable gamesrv@minecraft-smp
# Start:   systemctl start  gamesrv@minecraft-smp
# Logs:    journalctl -u gamesrv@minecraft-smp -f
[Unit]
Description=Game server: %i
After=network-online.target
Wants=network-online.target
# Manager may pull /etc/gamesrv.env for GAMESRV_INSTALL_ROOT overrides
# but per-server env lives in the server's own install dir.

[Service]
Type=simple
User=gamesrv
Group=gamesrv
WorkingDirectory=/srv/gameservers/%i
EnvironmentFile=-/srv/gameservers/%i/server.env
ExecStart=/srv/gameservers/%i/start.sh
# stop.sh must be graceful (send console 'stop' via tmux, wait). Systemd's own
# SIGTERM would corrupt Minecraft worlds — the stop script MUST wait for a clean exit.
ExecStop=/srv/gameservers/%i/stop.sh
# Modded packs (Forge with Deeper Darker, Create, etc.) can take 2+ minutes
# to flush custom dimensions during 'stop'. Give the whole stop sequence
# 5 minutes before systemd escalates to SIGTERM/SIGKILL — otherwise mid-save
# world corruption is very real. This value must be >= (stop.sh worst-case
# wait + a small buffer for the JVM to actually exit after tmux dies). See
# app/control.py::_STOP_TIMEOUT_SEC which polls systemctl for up to 240s.
TimeoutStopSec=300
Restart=on-failure
RestartSec=10
KillMode=mixed
# Per-server RAM backstop — per §2 the game must have a Java -Xmx or similar
# already; MemoryMax is the *last* line of defense so one server can't kill
# its neighbors. Tune via the manager (memory_mb) then re-run the template
# regeneration (bootstrap ensure-units).
MemoryMax=8G

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
