# systemd/gamesrv@.service
# /etc/systemd/system/gamesrv@.service
# Template unit: one instance per game server.  Instance name (%i) matches the
# server's `name:` in servers/<name>.yml AND its install dir under /srv/gameservers.
#
# Enable:  systemctl enable gamesrv@minecraft-smp
# Start:   systemctl start  gamesrv@minecraft-smp
# Logs:    journalctl -u gamesrv@minecraft-smp -f
[Unit]
Description=Game server: %i
After=network-online.target
Wants=network-online.target
# Manager may pull /etc/gamesrv.env for GAMESRV_INSTALL_ROOT overrides
# but per-server env lives in the server's own install dir.

[Service]
Type=simple
# Unreal Engine dedicated servers (Palworld, ARK) open tens of thousands of
# file handles while streaming world assets during load. The default 1024
# soft limit makes world-load STALL silently — RSS climbs to ~800 MB-1 GB
# then freezes with no crash and no error (the exact Palworld symptom).
# Every working native-Linux Palworld deploy raises this. start.sh also runs
# `ulimit -n` as a backstop in case this unit file isn't redeployed.
LimitNOFILE=100000
User=gamesrv
Group=gamesrv
WorkingDirectory=/srv/gameservers/%i
EnvironmentFile=-/srv/gameservers/%i/server.env
ExecStart=/srv/gameservers/%i/start.sh
# stop.sh must be graceful (send console 'stop' via tmux, wait). Systemd's own
# SIGTERM would corrupt Minecraft worlds — the stop script MUST wait for a clean exit.
ExecStop=/srv/gameservers/%i/stop.sh
# Modded packs (Forge with Deeper Darker, Create, etc.) can take 2+ minutes
# to flush custom dimensions during 'stop'. Give the whole stop sequence
# 5 minutes before systemd escalates to SIGTERM/SIGKILL — otherwise mid-save
# world corruption is very real. This value must be >= (stop.sh worst-case
# wait + a small buffer for the JVM to actually exit after tmux dies). See
# app/control.py::_STOP_TIMEOUT_SEC which polls systemctl for up to 240s.
TimeoutStopSec=300
# Restart=no on purpose. `Restart=on-failure` was masking every game-side
# crash by silently re-executing start.sh every RestartSec seconds — the
# operator saw "◐ starting" (Uptime 0s, RAM low) forever because
# ActiveEnterTimestamp reset on every restart, but no message ever
# appeared in the manager journal because these restarts happen INSIDE
# systemd, not through control.stop(). With Restart=no, a game that dies
# stays dead: the systemd unit flips to "failed" (visible as ✕ failed on
# the dashboard card) and start.sh's exit-code capture in the game's
# console.log gives the actual reason. Operators can then click Start
# again explicitly. If a specific server needs auto-restart, add a
# per-server drop-in:
#   mkdir -p /etc/systemd/system/gamesrv@<name>.service.d
#   printf '[Service]\nRestart=on-failure\nRestartSec=10\n' > \
#     /etc/systemd/system/gamesrv@<name>.service.d/restart.conf
#   systemctl daemon-reload
Restart=no
RestartSec=10
KillMode=mixed
# NO hardcoded MemoryMax here on purpose. A blanket cgroup memory cap in
# the template ignored the per-server memory_mb field entirely and hard-
# killed larger games (Palworld 16 GB, ARK 20+ GB) with an OOM signal
# the moment world-load exceeded the template default. The manager
# instead budgets RAM via a pre-flight check in app/main.py::api_action
# ("start" / "restart"): if the sum of active servers' memory_mb + the
# one being started would exceed the LXC's cgroup limit, the API refuses
# with 507 BEFORE systemctl is called. That protects the LXC from
# overcommit without silently killing individual games mid-run.
#
# If you want a per-server hard cap as belt-and-suspenders backstop,
# drop a systemd override alongside this template (Proxmox host, as
# root, adjust for your CT):
#   mkdir -p /etc/systemd/system/gamesrv@palworld.service.d
#   printf '[Service]\nMemoryMax=18G\n' > \
#     /etc/systemd/system/gamesrv@palworld.service.d/memory.conf
#   systemctl daemon-reload

# Hardening
NoNewPrivileges=true
# PrivateTmp intentionally NOT set. Steam Runtime (used by any SteamCMD
# game that links against libsteam_api.so — Palworld, Satisfactory, ARK,
# Valheim, Enshrouded) writes crash dumps, IPC pipes, and Steam runtime
# state into /tmp/dumps/ and /tmp/steam_pipe_*, and its SteamAPI_Init()
# expects those paths to be shared with other Steam processes on the
# same host. With PrivateTmp=true, systemd gives the unit a private /tmp
# namespace which silently breaks Steam networking init:
#   [S_API FAIL] Tried to access Steam interface SteamUser021 before SteamAPI_Init succeeded.
#   [S_API FAIL] Tried to access Steam interface SteamFriends017 ...
#   [S_API FAIL] Tried to access Steam interface STEAMAPPS_INTERFACE_VERSION008 ...
#   [S_API FAIL] Tried to access Steam interface SteamNetworkingUtils004 ...
# The game still binds its game port and prints "Running <game> on :PORT",
# but its Steam GameNetworkingSockets layer never comes up — so incoming
# client connections are silently dropped at the game layer and A2S
# queries queue on 27015 with a growing Recv-Q because Steam's built-in
# A2S handler isn't running. Result: dashboard stays "◐ starting" forever
# and clients see "Connection timed out". Removing PrivateTmp costs
# nothing (games don't leak sensitive data via /tmp) and unblocks every
# SteamCMD game on this template.

[Install]
WantedBy=multi-user.target
