# systemd/gamesrv-manager.service
# /etc/systemd/system/gamesrv-manager.service
# The FastAPI manager. Runs as unprivileged `gamesrv`. Env from /etc/gamesrv.env.
[Unit]
Description=Game Server Manager (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gamesrv
Group=gamesrv
EnvironmentFile=/etc/gamesrv.env
WorkingDirectory=/opt/gamesrv
ExecStart=/opt/gamesrv/run_manager.sh
# Restart=always so the Admin → Restart button (which just calls os._exit(0))
# actually brings the manager back. Also survives crashes.
Restart=always
RestartSec=5

# Hardening — Preston runs this LAN-only behind UFW, but defense in depth.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
# ProtectSystem=full mounts /etc read-only for us. Punch a hole for the
# one file the Admin env editor needs to write to.
ReadWritePaths=/etc/gamesrv.env

[Install]
WantedBy=multi-user.target
